# rootfs/etc/services.d/openwrt2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

bashio::log.info 'Link data dir'
rm -fr /opt/openwrt2mqtt/storage
ln -s /data /opt/openwrt2mqtt/storage

cd /opt/openwrt2mqtt

cat <<EOT >> config/settings.local.yml
openwrt:
  url: '$(bashio::config 'openwrt_url')'
  username: '$(bashio::config 'openwrt_username')'
  password: '$(bashio::config 'openwrt_password')'
  latitude: '$(bashio::config 'openwrt_latitude')'
  longitude: '$(bashio::config 'openwrt_longitude')'
mqtt:
  host: '$(bashio::config 'mqtt_host')'
  port: $(bashio::config 'mqtt_port')
  username: '$(bashio::config 'mqtt_username')'
  password: '$(bashio::config 'mqtt_password')'
  retain: $(bashio::config 'mqtt_retain')
debug: $(bashio::config 'debug')
EOT


if [ ! -f storage/production.sqlite3 ]; then
    bashio::log.info 'Create OpenWRT 2 MQTT database ...'
    bundle exec rails db:create
fi

bashio::log.info 'Migrate OpenWRT 2 MQTT database ...'
bundle exec rails db:migrate

bashio::log.info 'Start OpenWRT 2 MQTT ...'
bundle exec rails server
