# rootfs/etc/services.d/openwrt2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare APP__OPENWRT__URL
declare APP__OPENWRT__USERNAME
declare APP__OPENWRT__PASSWORD

declare APP__MQTT__HOST
declare APP__MQTT__PORT
declare APP__MQTT__USERNAME
declare APP__MQTT__PASSWORD
declare APP__MQTT__RETAIN

declare APP__DEBUG

## Get the 'message' key from the user config options.
export APP__OPENWRT__URL=$(bashio::config 'openwrt_url')
export APP__OPENWRT__USERNAME=$(bashio::config 'openwrt_username')
export APP__OPENWRT__USERNAME=$(bashio::config 'openwrt_password')

export APP__MQTT__HOST=$(bashio::config 'mqtt_host')
export APP__MQTT__PORT=$(bashio::config 'mqtt_port')
export APP__MQTT__USERNAME=$(bashio::config 'mqtt_username')
export APP__MQTT__PASSWORD=$(bashio::config 'mqtt_password')
export APP__MQTT__RETAIN=$(bashio::config 'mqtt_retain')

export APP__DEBUG=$(bashio::config 'debug')

bashio::log.info 'Link data dir'
rm -fr /opt/openwrt2mqtt/storage
ln -s /data /opt/openwrt2mqtt/storage

cd /opt/openwrt2mqtt

if [ ! -f storage/production.sqlite3 ]; then
    bashio::log.info 'Create OpenWRT 2 MQTT database ...'
    bin/rails db:create
fi

bashio::log.info 'Migrate OpenWRT 2 MQTT database ...'
bin/rails db:migrate

bashio::log.info 'Start OpenWRT 2 MQTT ...'
bin/rails server
